# scripts/launchd/com.nebula.watchpush.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.nebula.watchpush</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>/Users/mr.hu/Desktop/剧本分析/scripts/watch_and_push.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/mr.hu/Desktop/剧本分析</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
      <key>SERVER</key>
      <string>root@60.205.161.135</string>
      <key>REMOTE_DIR</key>
      <string>/opt/nebula</string>
      <key>PORT</key>
      <string>5000</string>
      <key>GIT_ENABLE</key>
      <string>true</string>
      <key>USE_LAUNCHD_WATCH</key>
      <string>true</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>
    <!-- WatchPaths 模式：目录变更触发执行一次脚本 -->
    <key>WatchPaths</key>
    <array>
      <string>/Users/mr.hu/Desktop/剧本分析</string>
    </array>
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/tmp/nebula_watch.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/nebula_watch.err</string>
  </dict>
</plist>